# deploy/systemd/idol-auth-backup.timer
[Unit]
Description=Run idol-auth PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
